# dist/ddnss-update.timer
[Unit]
Description=ddnss automatic update

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=basic.target
