# dist/ddnss-update.service
[Unit]
Description=ddnss automatic update

[Service]
User=ddnss
Type=oneshot
ExecStart=/usr/bin/ddnss-update
